# dotfiles/files/default/profile.d/pythonrc.sh
[ -f ~/.pythonrc.py ] && PYTHONSTARTUP=~/.pythonrc.py